# roles/storage/files/autorip/etc/systemd/system/autorip.service
[Unit]
Description=Autorip CD

[Service]
Type=oneshot
ExecStart=abcde /dev/cdrom -N -p -o flac
