# pi-image/cce-firstboot.service
[Unit]
Description=First-boot setup

Requires=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/opt/firstboot.sh

