# roles/storage-client/files/mnt-bigdata-serving.automount
[Unit]
Description=samba automount for bigdata/serving

[Automount]
Where=/mnt/bigdata/serving
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
