# roles/storage-client/files/mnt-bigdata-perpetual.automount
[Unit]
Description=samba automount for bigdata/perpetual
After=network-online.target
Wants=network-online.target
Requires=tailscaled.service
After=tailscaled.service

[Automount]
Where=/mnt/bigdata/perpetual
TimeoutIdleSec=5m

[Install]
WantedBy=multi-user.target
