# roles/storage-client/files/mnt-bigdata-serving.automount
[Unit]
Description=samba automount for bigdata/serving
After=network-online.target
Wants=network-online.target
Requires=tailscaled.service
After=tailscaled.service


[Automount]
Where=/mnt/bigdata/serving
TimeoutIdleSec=5m

[Install]
WantedBy=multi-user.target
