# roles/storage-client/files/mnt-bigdata-perpetual.automount
[Unit]
Description=samba automount for bigdata/perpetual

[Automount]
Where=/mnt/bigdata/perpetual
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
