# scripts/pi-image/cce-firstboot.service
[Unit]
Description=First-boot setup

WantedBy=multi-user.target
Requires=network-online.target

[Service]
ExecStart=/opt/firstboot.sh

